# extra/fan_pwm_100.timer
[Unit]
Description=Set speed fan to 100%


[Timer]
OnCalendar=08:00
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target